# devops/docker/php-with-code/supervisor-example.ini
;默认队列及异步事件
[program:laravel-queue-default]
command=php /opt/src/artisan queue:listen --delay=0 --memory=128 --sleep=30 --tries=3 --timeout=120
autostart=true
autorestart=true
user=root
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/crontab/module-default.log
